# config/launchd/com.astock.trade.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.astock_trading.trade.morning</string>
      <key>ProgramArguments</key>
  <array>
    <string>/Users/USER/.local/bin/atrade</string>
    <string>run-pipeline</string>
    <string>morning</string>
    <string>--json</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>25</integer>
    <key>Weekday</key><integer>1</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/USER/.local/state/a-stock-trading/logs/launchd/morning.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/USER/.local/state/a-stock-trading/logs/launchd/morning.err.log</string>
</dict>
</plist>
