# config/launchd/com.astock.trade.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.astock_trading.trade.health</string>
      <key>ProgramArguments</key>
  <array>
    <string>/Users/USER/.local/bin/atrade</string>
    <string>health</string>
    <string>--json</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StandardOutPath</key>
  <string>/Users/USER/.local/state/a-stock-trading/logs/launchd/health.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/USER/.local/state/a-stock-trading/logs/launchd/health.err.log</string>
</dict>
</plist>
